# autostart/redshift.desktop
[Desktop Entry]
Type=Application
Name=Redshift
Exec=redshift

