# runit/service/runsvdir-darrell/run
#!/bin/sh

export USER="darrell"
export HOME="/home/$USER"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/.config/runit_services"

exec 2>&1
exec chpst -u "$USER:$groups" runsvdir "$svdir"
